# deploy/opspanel.service
[Unit]
Description=OpsPanel — server management console
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Change these to match your install path / user.
WorkingDirectory=/opt/opspanel
ExecStart=/usr/bin/node --experimental-sqlite server/index.js
Restart=always
RestartSec=5
# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
Environment=NODE_ENV=production
# Optional: bind only to localhost and front with nginx instead.
# Environment=server.host=127.0.0.1

[Install]
WantedBy=multi-user.target
